# supervisord.conf
; Single-container process manager.
; Runs the main API and the Vite frontend as INDEPENDENT programs. Each has
; autorestart=true, so a crash restarts ONLY that program — the other keeps
; running.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

; Enable `supervisorctl` inside the container (restart individual programs).
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:api]
command=/opt/venv/bin/python -m uvicorn backend.api:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frontend]
command=npm run dev -- --host 0.0.0.0 --port 5173
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
